# rpcooler/usr/lib/systemd/system/rpcooler.service
[Unit]
Description=Cooler for Raspberry Pi

[Service]
ExecStart=/bin/rpcooler
Restart=always
RestartSec=5s

[Install]
WantedBy=multi-user.target
